# nginx/nginx-crl.timer
[Unit]
Description=Renewing CRL weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=multi-user.target
